# openlit-controller/deploy/openlit-controller.service
[Unit]
Description=OpenLIT Controller - LLM Observability via eBPF
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/openlit-controller --config /etc/openlit-controller/config.yaml
EnvironmentFile=-/etc/openlit-controller/env
Restart=always
RestartSec=5
AmbientCapabilities=CAP_BPF CAP_SYS_PTRACE CAP_NET_RAW CAP_PERFMON CAP_DAC_READ_SEARCH CAP_CHECKPOINT_RESTORE

[Install]
WantedBy=multi-user.target
